# web/where-server.service
[Unit]
Description=Where Tracking Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/where-web
EnvironmentFile=-/opt/where-web/.env
ExecStart=/usr/local/bin/bun --watch src/server/index.ts
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/where-web/data

[Install]
WantedBy=multi-user.target

